# ansible/roles/automount_ericg/files/home-ericg.automount
[Unit]
Description=home automount for ericg

[Automount]
Where=/home/ericg

[Install]
WantedBy=multi-user.target
